# build/pacmon.service
[Unit]
Description=pacmon

[Service]
ExecStart=/usr/bin/pacmon
User=pacmon
Restart=always

[Install]
WantedBy=multi-user.target
